# pwrstat-api.service
[Unit]
Description=pwrstat REST API
Documentation=https://github.com/chattersley/pwrstat-python-api
After=network.target pwrstatd.service
Wants=pwrstatd.service

[Service]
Type=simple
# Runs as root: pwrstatd only permits root access to its IPC socket.
User=root
Group=root

ExecStart=/opt/pwrstat-api/.venv/bin/pwrstat-api
Restart=on-failure
RestartSec=5

# Configuration — override in /etc/systemd/system/pwrstat-api.service.d/local.conf
Environment=PORT=5002
Environment=HOST=0.0.0.0
Environment=LOG_LEVEL=info

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
